# gui-client/systemd/bluebubbles-preview-proxy.service
[Unit]
Description=BlueBubbles Preview Proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Code/bluebubbles-gui
ExecStart=%h/Code/bluebubbles-gui/bluebubbles-preview-proxy
Restart=on-failure
RestartSec=3
Environment=PREVIEW_PROXY_ADDR=127.0.0.1:8090
Environment=PREVIEW_OEMBED_ENDPOINT=https://noembed.com/embed
Environment=PREVIEW_TIMEOUT_SEC=8
Environment=PREVIEW_CACHE_TTL_SEC=21600

[Install]
WantedBy=default.target
